# config/com.evrone.masquito.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.evrone.masquito</string>
  <key>ProgramArguments</key>
  <array>
    <string><%= RbConfig.ruby %></string>
    <string><%= ::Masquito::BIN %></string>
    <string>start</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CONFIG_PATH</key>
    <string><%= ::Masquito::CONFIG_PATH %></string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
